# Dockerfiles/rootfs/etc/crontabs/db-backup
0 4 * * * /etc/service.d/mariadb/bin/backup > /dev/null 2>&1
