# 10.2/rootfs/etc/crontabs/db-backup
0 4 * * * /etc/service.d/mariadb/backup > /dev/null 2>&1
